# bin/templates/etc/rc.d/freedive
#!/bin/sh
# shellcheck disable=SC1091,2034

# PROVIDE: freedive
# REQUIRE: LOGIN
#
# Add the following line to /etc/rc.conf to enable freedive:
#
# freedive_enable="YES"
#

freedive_enable=${freedive_enable:-"NO"}
. /etc/rc.subr

name=freedive
rcvar=freedive_enable

extra_commands="run iex init config import export reset" 

procname=/usr/local/bin/freedive
start_cmd="/usr/local/bin/freedive daemon"
stop_cmd="/usr/local/bin/freedive stop"

run_cmd="/usr/local/bin/freedive start"
iex_cmd="/usr/local/bin/freedive start_iex"

init_cmd="/usr/local/freedive/scripts/init.sh"
config_cmd="/usr/local/freedive/scripts/config.sh"
import_cmd="/usr/local/freedive/scripts/import.sh"
export_cmd="/usr/local/freedive/scripts/export.sh"
reset_cmd="/usr/local/freedive/scripts/reset.sh"

load_rc_config ${name}
run_rc_command "$1"
